# bin/wetty.service
[Unit]
Description=Wetty Web Terminal
After=network.target

[Service]
ExecStart=/root/wetty/index.js -p 3000 --host 127.0.0.1 -b "" --forcessh --title ""

[Install]
WantedBy=multi-user.target
